# tools/com.voicebadge.weather.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.voicebadge.weather</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/taohe/work/vbctrl/tools/weather-venv/bin/python</string>
        <string>/Users/taohe/work/vbctrl/tools/mac_weather_push.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>StandardOutPath</key>
    <string>/tmp/voicebadge-weather.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/voicebadge-weather.log</string>
</dict>
</plist>
